# roles/kubesolo/templates/kubesolo.openrc.j2
#!/sbin/openrc-run
# KubeSolo OpenRC service script
# Generated by Ansible

name="kubesolo"
description="KubeSolo - Lightweight Kubernetes Distribution"

command="{{ kubesolo_install_path }}"
command_args="--path={{ kubesolo_path }}{% if kubesolo_apiserver_extra_sans | length > 0 %} --apiserver-extra-sans={{ kubesolo_apiserver_extra_sans | join(',') }}{% endif %}{% if kubesolo_portainer_edge_id %} --portainer-edge-id={{ kubesolo_portainer_edge_id }}{% endif %}{% if kubesolo_portainer_edge_key %} --portainer-edge-key={{ kubesolo_portainer_edge_key }}{% endif %}{% if kubesolo_portainer_edge_async %} --portainer-edge-async{% endif %}{% if kubesolo_local_storage %} --local-storage{% endif %}{% if kubesolo_debug %} --debug{% endif %}{% if kubesolo_pprof_server %} --pprof-server{% endif %}"

command_background=true
pidfile="{{ kubesolo_pidfile }}"
output_log="{{ kubesolo_logfile }}"
error_log="{{ kubesolo_logfile }}"

# Service dependencies
depend() {
    need net
    after firewall
}

# Pre-start checks
start_pre() {
    # Ensure log directory exists
    mkdir -p "$(dirname {{ kubesolo_logfile }})"

    # Ensure PID directory exists
    mkdir -p "$(dirname {{ kubesolo_pidfile }})"

    # Set ulimits
    ulimit -n 65535

    return 0
}

# Environment variables
{% if kubesolo_proxy %}
export HTTP_PROXY="{{ kubesolo_proxy }}"
export HTTPS_PROXY="{{ kubesolo_proxy }}"
export NO_PROXY="localhost,127.0.0.1,0.0.0.0,10.0.0.0/8,192.168.0.0/16"
{% endif %}

# Service behavior
respawn_delay=5
respawn_max=0
retry="TERM/30/KILL/5"

# OOM adjustment
oom_score_adj=-999
