# roles/kubesolo/templates/kubesolo.upstart.j2
# KubeSolo Upstart configuration
# Generated by Ansible
# Place this file at /etc/init/kubesolo.conf

description "KubeSolo - Lightweight Kubernetes Distribution"
author "Portainer.io"

# Start/stop conditions
start on runlevel [2345]
stop on runlevel [!2345]

# Restart behavior
respawn
respawn limit 10 5

# Process management
kill timeout 30

# Set resource limits
limit nofile 65535 65535
limit nproc unlimited unlimited

# Environment variables
{% if kubesolo_proxy %}
env HTTP_PROXY="{{ kubesolo_proxy }}"
env HTTPS_PROXY="{{ kubesolo_proxy }}"
env NO_PROXY="localhost,127.0.0.1,0.0.0.0,10.0.0.0/8,192.168.0.0/16"
{% endif %}

# Pre-start script
pre-start script
    # Ensure directories exist
    mkdir -p "$(dirname {{ kubesolo_logfile }})"
    mkdir -p "$(dirname {{ kubesolo_pidfile }})"
end script

# Main script
script
    # Build command arguments
    CMD_ARGS="--path={{ kubesolo_path }}"
{% if kubesolo_apiserver_extra_sans | length > 0 %}
    CMD_ARGS="$CMD_ARGS --apiserver-extra-sans={{ kubesolo_apiserver_extra_sans | join(',') }}"
{% endif %}
{% if kubesolo_portainer_edge_id %}
    CMD_ARGS="$CMD_ARGS --portainer-edge-id={{ kubesolo_portainer_edge_id }}"
{% endif %}
{% if kubesolo_portainer_edge_key %}
    CMD_ARGS="$CMD_ARGS --portainer-edge-key={{ kubesolo_portainer_edge_key }}"
{% endif %}
{% if kubesolo_portainer_edge_async %}
    CMD_ARGS="$CMD_ARGS --portainer-edge-async={{ kubesolo_portainer_edge_async }}"
{% endif %}
{% if kubesolo_local_storage %}
    CMD_ARGS="$CMD_ARGS --local-storage=true"
{% endif %}
{% if kubesolo_debug %}
    CMD_ARGS="$CMD_ARGS --debug=true"
{% endif %}
{% if kubesolo_pprof_server %}
    CMD_ARGS="$CMD_ARGS --pprof-server=true"
{% endif %}

    # Execute kubesolo
    exec {{ kubesolo_install_path }} $CMD_ARGS >> {{ kubesolo_logfile }} 2>&1
end script

# Post-stop script
post-stop script
    # Clean up PID file if it exists
    rm -f {{ kubesolo_pidfile }}
end script
