# roles/kubesolo/templates/kubesolo.service.j2
# KubeSolo systemd service unit file
# Generated by Ansible

[Unit]
Description=KubeSolo - Lightweight Kubernetes Distribution
Documentation=https://github.com/portainer/kubesolo
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart={{ kubesolo_install_path }} --path={{ kubesolo_path }}{% if kubesolo_apiserver_extra_sans | length > 0 %} --apiserver-extra-sans={{ kubesolo_apiserver_extra_sans | join(',') }}{% endif %}{% if kubesolo_portainer_edge_id %} --portainer-edge-id={{ kubesolo_portainer_edge_id }}{% endif %}{% if kubesolo_portainer_edge_key %} --portainer-edge-key={{ kubesolo_portainer_edge_key }}{% endif %}{% if kubesolo_portainer_edge_async %} --portainer-edge-async{% endif %}{% if kubesolo_local_storage %} --local-storage{% endif %}{% if kubesolo_debug %} --debug{% endif %}{% if kubesolo_pprof_server %} --pprof-server{% endif %}

KillMode=process
Restart=always
RestartSec=5s
LimitNOFILE=65535
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
TimeoutStartSec=0
Delegate=yes
{% if kubesolo_proxy %}
Environment="HTTP_PROXY={{ kubesolo_proxy }}"
Environment="HTTPS_PROXY={{ kubesolo_proxy }}"
Environment="NO_PROXY=localhost,127.0.0.1,0.0.0.0,10.0.0.0/8,192.168.0.0/16"
{% endif %}

# Security settings
# Run with reduced OOM killer priority
OOMScoreAdjust=-999

[Install]
WantedBy=multi-user.target
